# libs/vnet/examples/v1.25/motlie-vfs-guest.service
[Unit]
Description=motlie-vfs guest filesystem mounter (v1.15 Vz)
After=local-fs.target cloud-final.service
Wants=cloud-final.service

[Service]
Type=simple
ExecStart=/usr/local/bin/motlie-vfs-guest-v1_15 /etc/motlie-vfs/mounts.yaml
Restart=on-failure
RestartSec=2
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
